# deploy/launchd/ai.rick.telegram-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.rick.telegram-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__RICK_WORKSPACE__/scripts/run-telegram-bridge.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__RICK_WORKSPACE__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>__RICK_DATA_ROOT__/logs/rick-telegram-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>__RICK_DATA_ROOT__/logs/rick-telegram-bridge.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RICK_ENV_FILE</key>
    <string>__RICK_WORKSPACE__/config/rick.env</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
